# setup/crontab
1 * * * * php /var/travian.idrinth.de/bin/world-import.php